# frontend/supervisord.conf
[supervisord]
nodaemon=true
pidfile=/tmp/supervisord.pid
logfile=/dev/null
logfile_maxbytes=0

[program:caddy]
command=caddy run --config /etc/caddy/Caddyfile --adapter caddyfile
autostart=true
autorestart=false
environment=HOME="/tmp"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nextjs]
command=node /app/start.js
autostart=true
autorestart=false
environment=PORT="3001",HOSTNAME="127.0.0.1"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# If either program exits, take down supervisord so K8s sees the Pod as
# unhealthy and restarts it. Without this, a Next.js crash would leave Caddy
# happily serving 502s.
[eventlistener:quit_on_exit]
command=/bin/sh -c "printf 'READY\n'; IFS= read -r _; kill -TERM $PPID; printf 'RESULT 2\nOK'"
events=PROCESS_STATE_EXITED,PROCESS_STATE_FATAL
stdout_logfile=/dev/null
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
